# server/resources/cron-jobs/portal_acct.cron
0 0 * * *   /usr/bin/python /mnt/corn/sumarize_radacct_daily.py >/dev/null 2>&1
# 0 0 30 * *  /usr/bin/python /mnt/corn/sumarize_radacct_monthly.py >/dev/null 2>&1
