# server/resources/cron-jobs/summarizer.cron
0 0 * * *   /usr/bin/python /mnt/corn/daily_summarizer.py &
# 0 0 30 * *  /usr/bin/python /mnt/corn/summarize_radacct_monthly.py >/dev/null 2>&1
